# .deploy/docker/crontab
* * * * * root /artisan schedule:run >> /var/log/cron.log
